# pictureframe.service
[Unit]
Description=Picture Frame - Digital Picture Frame Application
After=network.target
Wants=network-online.target

[Service]
Type=simple
# Port 80 erfordert root-Rechte
# Lauf als root für vollständige Funktionalität
User=root
Group=root
WorkingDirectory=/home/pi/PictureFrameV3
Environment="DISPLAY=:0"
Environment="XAUTHORITY=/var/run/lightdm/root/:0"
Environment="XDG_RUNTIME_DIR=/run/user/0"
# Warte auf LightDM/X11
ExecStartPre=/bin/sleep 10
# Verwende setcap-enabled Python für Port 80
ExecStart=/usr/bin/python3 /home/pi/PictureFrameV3/src/main.py
Restart=always
RestartSec=10
# Starte nach graphical.target (wenn verfügbar) oder network.target
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

