# x11.service
[Unit]
Description=X11 Server for PictureFrame
After=network.target
Before=pictureframe.service

[Service]
Type=simple
User=root
ExecStart=/usr/bin/X :0 -nolisten tcp vt7
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target

